# linux/plug.service
[Unit]
Description=Balanced Mode on Plug

[Service]
User=root
ExecStart=/home/maneesh/.config/scripts/powersave 1
